# factory/pmaports/postmarketos-boot-factorytest/postmarketos-boot-factorytest.post-install
#!/bin/sh

rc-update add lightdm default
